# config/systemd/user/live-notifyd.service
[Unit]
Description=Live Notify daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/bin/live-notifyd serve
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
